# systemd/beauties-cleanup.timer
[Unit]
Description=Beauties files cleanup

[Timer]
OnCalendar=*:0/30
Persistent=true
Unit=beauties-cleanup.service

[Install]
WantedBy=timers.target